# usr/lib/opkg/info/nvram.list
/etc/init.d/nvram
/usr/lib/libnvram.so.0.1
/usr/sbin/nvram
