# usr/lib/opkg/info/vnstat.list
/etc/init.d/vnstat
/usr/bin/vnstat
/etc/config/vnstat
/etc/vnstat.conf
/usr/sbin/vnstatd
/etc/uci-defaults/vnstat
